# etc/update-motd.d/02-purrenjin
#!/bin/bash

# 可选猫脸（随机出现）
CAT_FACES=(
    "(=^･ω･^=)"
    "(ΦωΦ)"
    "(=①ω①=)"
    "(=ＴェＴ=)"
    "(=ↀωↀ=)"
    "(ฅ^•ﻌ•^ฅ)"
    "(=ಠωಠ=)"
    "(=ʘᆽʘ=)"
    "(=^◡^=)"
    "(=｀ω´=)"
)

# 随机选一个猫脸
RANDOM_FACE=${CAT_FACES[$RANDOM % ${#CAT_FACES[@]}]}

# 猫图（不含呼噜）
CAT_ASCII="
   ))    wWw  wWw()_() ()_()       wWw \\\  ///   _ wW  Ww\\\  ///
  (o0)-. (O)  (O)(O o) (O o)       (O)_((O)(O)) _||\(O)(O)((O)(O))
   | (_))/ )  ( \ |^_\  |^_\       / __)| \ || (_'\  (..)  | \ || 
   | .-'/ /    \ \|(_)) |(_))     / (   ||\\||  (  |  ||   ||\\|| 
   |(   | \____/ ||  /  |  /     (  _)  || \ |   \ | _||_  || \ | 
    \)  '. \`--' .\`)|\\  )|\\      \ \_  ||  ||(\__)|(_/\_) ||  || 
    (     \`-..-' (/  \)(/  \)      \__)(_/  \_)--.)      (_/  \_)
         $RANDOM_FACE   — PurrEnjin
         Morse: .--. ..- .-. .-.    . -. .--- .. -.
"

# 呼噜波纹数组（呼吸节奏）
PURR_WAVES=(
    "        .oOo. purr .oOo."
    "         oOo. purr .oOo "
    "          Oo. purr .oO  "
    "         oOo. purr .oOo "
)

# 呼噜动画
purr_breathe() {
    for cycle in {1..2}; do
        for wave in "${PURR_WAVES[@]}"; do
            clear
            echo "$CAT_ASCII" | lolcat -a -d 6 -s 80
            echo -e "\n$wave" | lolcat -a -d 6 -s 80
            sleep 0.15
        done
    done
}

# 播放动画
purr_breathe

# 显示系统信息
neofetch --ascii_distro "Arch" --ascii_colors 4 5 \
         --disable packages resolution \
         uptime shell cpu memory | lolcat -a -d 5 -s 80
